# templates/default/rc.local.erb
#!/bin/sh -e
#
# rc.local - Managed by Chef #
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits. It should be noted that rc.local has been replaced by systemd
# in future operating systems (Debian 9 for example)
#
#
/home/<%=node[:team][:admin]%>/miner_launcher.sh 15 &
exit 0
